# files/default/iptables.upstart.conf
description "IP Tables Firewall"

start on runlevel [2345]
stop on runlevel [016]
console log

task
exec /sbin/iptables-restore < /etc/iptables.rules